# supervisor/ssl.conf
;SSH service
[program:ssh]
directory=/
command=/usr/sbin/sshd -D
user=root
startretries=3
autorestart=true
redirect_stderr=true
autostart=true
stopasgroup=true
killasgroup=true
stopsignal=QUIT
